# vnstat-server.service
[Unit]
Description=vnstat HTTP Server
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/vnstat-http-server -port 8080 -token YOUR_TOKEN_HERE
# Optional: Enable Grafana Cloud push
# ExecStart=/usr/local/bin/vnstat-http-server -port 8080 -token YOUR_TOKEN_HERE -grafana-url "https://YOUR_PROMETHEUS_INSTANCE.grafana.net/api/prom/push" -grafana-user "YOUR_INSTANCE_ID" -grafana-token "YOUR_API_TOKEN" -grafana-interval 30s
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

